# service/piframe.service
[Unit]
Description=PiFrame Digital Photo Frame
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/piframe
ExecStart=/home/pi/piframe/venv/bin/python run.py

# Required for pygame to access the framebuffer without X11
Environment=SDL_VIDEODRIVER=kmsdrm
Environment=SDL_FBDEV=/dev/fb0

Restart=on-failure
RestartSec=10

# Allow access to GPIO / SPI for e-ink display
SupplementaryGroups=gpio spi i2c video

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
