# service/piframe.service
[Unit]
Description=PiFrame Digital Photo Frame
After=network-online.target graphical.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/PiFrame
ExecStart=/home/pi/PiFrame/venv/bin/python /home/pi/PiFrame/run.py

# ── Display environment ────────────────────────────────────────────────────────
# setup.sh patches XDG_RUNTIME_DIR with the correct UID automatically.

# Wayland desktop session (Pi OS Bookworm — default):
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/1000

# X11 desktop session (Pi OS Bullseye):
#Environment=DISPLAY=:0
#Environment=XAUTHORITY=/home/pi/.Xauthority

Restart=on-failure
RestartSec=5

# Allow access to GPIO / SPI for e-ink display
SupplementaryGroups=gpio spi i2c video

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
